# systemd/user/dwm.target
[Unit]
BindsTo=graphical-session.target
Requires=dwm.service
