# systemd/user/xbanish.service
[Install]
WantedBy=graphical-session.target

[Service]
Environment="XAUTHORITY=/run/user/1000/xauthority"
ExecStart=/home/matt/.nix-profile/bin/xbanish
Restart=on-failure
RestartSec=3

[Unit]
After=graphical-session-pre.target
Description=xbanish
PartOf=graphical-session.target
